# systemd/.config/systemd/user/gitsync.timer
[Unit]
Description=Trigger Git Sync once a day

[Timer]
# Run exactly once per day
OnCalendar=daily
# If the computer was off at midnight, run it as soon as you boot
Persistent=true
# Wait 2 minutes after you log in to let Wi-Fi connect
OnStartupSec=2min

[Install]
WantedBy=timers.target
